# netbird/rootfs/etc/s6-overlay/s6-rc.d/netbird/run
#!/command/with-contenv bashio
# shellcheck shell=bash

# 等待 HAOS network/system 完整啟動 v0.59.12
bashio::log.info "~~!! back to >> v0.59.12 Waiting extra 5 seconds for HAOS to stabilize..."
sleep 5

bashio::log.info "NetBird Client – Jason Custom Startup Script"

# ===============================================================
# 1. Wait for HAOS network to be ready
# ===============================================================
bashio::log.info "Waiting for network interface (max 30s)..."
for i in $(seq 1 30); do
    if ip route | grep default >/dev/null 2>&1; then
        bashio::log.info "Network interface ready."
        break
    fi
    sleep 1
done

# ===============================================================
# 2. Wait for DNS to resolve correctly
# ===============================================================
bashio::log.info "Checking DNS resolution (max 30s)..."
for i in $(seq 1 30); do
    if ping -c1 -W1 google.com >/dev/null 2>&1; then
        bashio::log.info "DNS OK."
        break
    fi
    sleep 1
done

# ===============================================================
# 3. Optional: Wait for NetBird Management server reachable
# ===============================================================
MANAGEMENT_URL="$(bashio::config 'management_url')"
if [ -n "$MANAGEMENT_URL" ]; then
    HOST=$(echo "$MANAGEMENT_URL" | awk -F[/:] '{print $4}')
    bashio::log.info "Checking NetBird management server: $HOST"

    for i in $(seq 1 20); do
        if ping -c1 -W1 "$HOST" >/dev/null 2>&1; then
            bashio::log.info "NetBird management server reachable."
            break
        fi
        sleep 1
    done
fi


# ==============================================================================
# Home Assistant Community Add-on: NetBird Client
# Runs NetBird Client
# ==============================================================================
declare -a options
declare name
declare value


# Get the options configured in HASS GUI
readonly CONFIG_OLD_PATH=/homeassistant/netbird/config.json
readonly CONFIG_PATH=/config/config.json

[ -f "${CONFIG_OLD_PATH}" ] && mv "${CONFIG_OLD_PATH}" "${CONFIG_PATH}"

admin_url="$(bashio::config 'admin_url')"
management_url="$(bashio::config 'management_url')"
setup_key="$(bashio::config 'setup_key')"
hostname="$(bashio::config 'hostname')"
rosenpass="$(bashio::config 'rosenpass')"
rosenpass_permissive="$(bashio::config 'rosenpass_permissive')"
log_level="$(bashio::config 'log_level')"

options+=(--foreground-mode)
options+=(--config "${CONFIG_PATH}")
options+=(--log-file console)

if [ "${admin_url}" = "" ]; then
    bashio::log.info "Using Default Admin URL"
else
    bashio::log.info "Using ${admin_url} as Admin URL"
    options+=(--admin-url "${admin_url}")
fi

if [ "${management_url}" = "" ]; then
    bashio::log.info "Using Default Management URL"
else
    bashio::log.info "Using ${management_url} as Management URL"
    options+=(--management-url "${management_url}")
fi

if [ "${setup_key}" = "" ]; then
    bashio::log.info "No Setup Key Set"
    bashio::log.info "This client will only show up in dashboards it's already registered with."
else
    bashio::log.info "Using ${setup_key} as Setup Key"
    options+=(--setup-key "${setup_key}")
fi

if [ "${hostname}" = "" ]; then
    bashio::log.info "No Hostname Set"
    bashio::log.info "This client will use the default (<docker container id>-netbird-client) as hostname in peers."
else
    bashio::log.info "Using ${hostname} as hostname"
    options+=(--hostname "${hostname}")
fi

if ! ${rosenpass}; then
    bashio::log.info "Rosenpass disabled"
    options+=(--enable-rosenpass=false)
else
    bashio::log.info "Rosenpass enabled"
    options+=(--enable-rosenpass)
    if ${rosenpass_permissive}; then
        bashio::log.info "Rosenpass permissive mode enabled"
        options+=(--rosenpass-permissive)
    fi
fi

if [ "${log_level}" = "" ] || [ "${log_level}" = "null" ]; then
    bashio::log.info "No log level Set"
    bashio::log.info "This client will use the default logging."
else
    bashio::log.info "Using ${log_level} as logging level"
    options+=(--log-level "${log_level}")
fi

# Load custom environment variables
for var in $(bashio::config 'env_vars|keys'); do
    name=$(bashio::config "env_vars[${var}].name")
    value=$(bashio::config "env_vars[${var}].value")
    bashio::log.info "Setting ${name} to ${value}"
    export "${name}=${value}"
done

# dirty hack to get dns working
# with this netbird can regonize the host running systemd-resolved
CONTENT=$(cat /etc/resolv.conf)
echo '# systemd-resolved' > /etc/resolv.conf
echo "$CONTENT" >> /etc/resolv.conf


# ---------------------------------------------------------------
# 顯示 Add-on 版本（config.yaml）
ADDON_VERSION=$(bashio::addon.version)
bashio::log.info "Add-on Version (config.yaml): ${ADDON_VERSION}"

# 顯示 Docker image 名稱
# DOCKER_IMAGE=$(bashio::addon.image)
# bashio::log.info "Docker Image: ${DOCKER_IMAGE}"
# ERROR 錯誤導致重啟 bashio::addon.image: command not found

# 顯示 NetBird 二進位版本
NETBIRD_BIN_VERSION=$(netbird version)
bashio::log.info "NetBird Binary Version: ${NETBIRD_BIN_VERSION}"

# 顯示 NetBird 介面 wt0 IP（若已建立）
# WT0_IP=$(ip -o -4 addr show wt0 2>/dev/null | awk '{print $4}')
# bashio::log.info "NetBird wt0 IP: ${WT0_IP:-not ready yet}"
# ERROR 會stop

# 顯示 netbird up 的參數
bashio::log.info "Starting NetBird Client..."
bashio::log.info "netbird up " "${options[@]}"

# 運行 netbird
netbird up "${options[@]}"
